# launchagents/com.sklavit.claude-usage.scan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sklavit.claude-usage.scan</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>$HOME/.claude-usage-archive/usage.py</string>
        <string>scan</string>
        <string>--export-dir</string>
        <string>$HOME/personal/personal_notes/_automation/claude-usage-tracker/data</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>9</integer>
        <key>Minute</key><integer>5</integer>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key><string>$HOME/.claude-usage-archive/scan.log</string>
    <key>StandardErrorPath</key><string>$HOME/.claude-usage-archive/scan.log</string>
</dict>
</plist>
